# conf/wrapper.service
[Unit]
Description=Wrapper

[Service]
Type=idle
Environment=NODE_PATH=/usr/local/lib/node_modules
ExecStartPre=/opt/central_monitor/script/gpu_env.sh
ExecStart=/usr/bin/nodejs /opt/central_monitor/wrapper/index.js 
ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGQUIT
TimeoutStopSec=5
KillMode=process
PrivateTmp=true
Restart=always
RestartSec=1s

[Install]
WantedBy=multi-user.target

